# templates/fetch-modules.service.erb
# This file is managed by Puppet
[Unit]
Description=Fetch Updates for Puppet Modules with g10k and Postrun
Requires=local-fs.target
ConditionPathIsReadWrite=/etc/puppetlabs/code/environments
ConditionFileIsExecutable=/usr/local/bin/g10k
ConditionFileNotEmpty=/etc/puppetlabs/r10k/postrun/postrun.py
ConditionFileNotEmpty=/etc/puppetlabs/r10k/g10k.yaml

[Service]
Type=oneshot
# deploy
ExecStart=/usr/local/bin/g10k -quiet -config /etc/puppetlabs/r10k/g10k.yaml
# postrun
ExecStart=/usr/bin/python3 /etc/puppetlabs/r10k/postrun/postrun.py


[Install]
WantedBy=multi-user.target
